# docs/deployment/agenthive-orchestrator.service
[Unit]
Description=AgentHive Orchestrator — Event-driven agent dispatcher
Documentation=https://github.com/agentroadmap/agentRoadmap
After=network-online.target postgresql.service agenthive-gate-pipeline.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/data/code/AgentHive
ExecStart=/usr/bin/env node --import jiti/register scripts/orchestrator.ts
Restart=on-failure
RestartSec=10
# P254: 5min stop window covers the orchestrator's 600s spawn timeout's
# tail so dispatchImplicitGate can finish before SIGKILL.
TimeoutStopSec=300
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agenthive-orchestrator

# Environment file for database credentials
EnvironmentFile=-%h/.agenthive.env

# P245: declare host identity so the spawner can look up host_model_policy
# and refuse forbidden routes (e.g. anthropic on hermes).
Environment=AGENTHIVE_HOST=hermes

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=default.target
