# scripts/systemd/agenthive-mcp.service
[Unit]
Description=AgentHive MCP SSE Server
After=network.target

[Service]
Type=simple
User=agenthive
Group=agenthive
WorkingDirectory=/data/code/AgentHive
EnvironmentFile=/etc/agenthive/env
Environment=HOME=/home/agenthive
Environment=NODE_ENV=production
Environment=PATH=/home/agenthive/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/agenthive-mcp.sh
Restart=always
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
